# launchagents/Library/LaunchAgents/com.kevin.voice-claude.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.kevin.voice-claude</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/kevin/.pyenv/versions/3.13.7/bin/python</string>
        <string>/Users/kevin/.ai/scripts/voice-claude/voice_menubar.py</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/kevin/.local/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
    <key>WorkingDirectory</key>
    <string>/Users/kevin/.ai</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/tmp/voice-claude.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/voice-claude.log</string>
</dict>
</plist>
